# handbandd.conf
description	"Hand Band printing daemon"

start on (started cups and started mysql)
stop on runlevel [!2345]

respawn

exec python /usr/share/handbandd/handbandd_und.py
